# files/systemd/com.webos.service.uwb.service.in
[Unit]
Description=webos - "%n"
Requires=ls-hubd.service
After=ls-hubd.service

[Service]
Type=simple
OOMScoreAdjust=-500
EnvironmentFile=-/var/systemd/system/env/com.webos.service.uwb.env
Environment=CHARSET=UTF-8
ExecStart=/usr/sbin/webos-uwb-service
#Restart=on-failure
